# config/systemd/jimbook.service
[Unit]
Description=JimBook

[Service]
WorkingDirectory=/opt/jimbook
ExecStart=/usr/bin/dotnet /opt/jimbook/JimBook.dll
Restart=always
RestartSec=10
SyslogIdentifier=JimBook
User=www-data
Environment=ASPNETCORE_ENVIRONMENT=Production 

[Install]
WantedBy=multi-user.target
